# MachineRoom/door
#!/bin/sh /etc/rc.common

START=90

SERVICE_DAEMONIZE=1

start(){
        service_start lua /usr/bin/serve.lua
}

stop() {
        service_stop lua /usr/bin/serve.lua
}
