# MachineRoom/ledoff
#!/bin/sh /etc/rc.common
# Turn off LED at startup prevent blinking

START=15
start() {
	# set leds to normal state
	. /etc/diag.sh
	set_state done
}

